# docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=warn
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
autostart=true
autorestart=unexpected
exitcodes=0
; php-fpm stops first (stop in reverse priority order)
priority=5
; SIGQUIT = graceful drain (finish in-flight requests)
; SIGTERM = immediate kill
stopsignal=QUIT
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=unexpected
exitcodes=0
; nginx stops second (higher priority value = stopped last in stop order)
priority=10
stopsignal=QUIT
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
; Laravel/Statamic task scheduler. Runs php artisan schedule:run every minute.
; Handles: session cleanup, content versioning hooks, any queued Statamic jobs.
; `schedule:work` is the long-running variant (replaces a system cron).
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=unexpected
; Scheduler starts last — nginx and php-fpm must be up first.
priority=15
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
